# apps/releases/etc/rc5.d/S01MuxLab.sh
#! /bin/sh
### BEGIN INIT INFO
### END INIT INFO


echo "Configure SFP interface..."
LOG_DIR=/var/log

sfpCtrl >> $LOG_DIR/sfp.log
